# misc/crontab
*/15 * * * * offlineimap; pkill -RTMIN+11 i3blocks
*/5 * * * * ~/.sync/sync.sh
@hourly sudo pacman -Sy; pkill -RTMIN+12 i3blocks
